# Suranta_Fly/main_server/flight-monitor-celery.service
[Unit]
Description=Flight Monitor Celery Worker and Beat
After=network.target postgresql.service redis.service

[Service]
User=flightmonitor
Group=flightmonitor
WorkingDirectory=/opt/flight-monitor
Environment="PATH=/opt/flight-monitor/venv/bin"
EnvironmentFile=/opt/flight-monitor/.env
ExecStart=/opt/flight-monitor/venv/bin/celery -A app.services.monitoring worker --loglevel=info --beat
Restart=always

[Install]
WantedBy=multi-user.target 